# nim_planet.timer
[Unit]
Description=Nim Planet

[Timer]
OnCalendar=*:0/10
